# runonce.service
[Unit]
Description=First Start (builds certs for sendmail)
Before=sendmail.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/runonce

[Install]
WantedBy=multi-user.target
